# profile.d/aliases.sh
# aliases.sh

# General
alias ..='cd ..'
alias ...='cd ../..'
alias ls='ls -FG'
alias ll='ls -lFG'
alias la='ls -aFG'
alias lla='ls -alFG'
alias mv='mv -i'
alias rm='rm -i'
alias du='du -h'
alias df='df -h'
alias tg='tar -zcf'
alias tb='tar -jcf'
alias tx='tar -Jcf'
alias utg='tar -zxf'
alias utb='tar -jxf'
alias utx='tar -Jxf'
