# profile.d/python.sh
# Python Settings
export PIP_PREFIX="$HOME/Library/Python/2.7"
if [ -r ~/.pythonrc.py ]; then
    export PYTHONSTARTUP=~/.pythonrc.py
fi
